# recipes-core/udev/demo-i2c-udev-conf/rtc-i2c-attach@.service
[Service]
Type=oneshot
ExecStartPre=/bin/echo '#i2c rtc device has been attached'
ExecStart=/usr/sbin/rtc-i2c-attach.sh %i
