# recipes-support/udisks/files/automount.service
[Unit]
Description=Automount Disk Manager
After=udisk.service

[Service]
ExecStart=/usr/libexec/automount.sh

[Install]
WantedBy=default.target
