# recipes-support/udisks/files/automount.service
[Unit]
Description=Automount Disk Manager
After=udisks.service

[Service]
ExecStart=/usr/libexec/automount.sh

[Install]
WantedBy=multi-user.target
